# dist/scripts/postinstall
#!/bin/sh
# Activate: fix kext ownership (kextload requires root:wheel), reset the boot
# sentinel to ok (a fresh install should attempt full-gesture mode), and
# start the daemon. The launchd wrapper (voodooinputmavericks-run) loads the one kext from
# /usr/local and re-enumerates the device so the in-kernel readers win their
# matches. No /Library/Extensions, no kextcache: the unsigned kext only loads
# from a loose path.
set -e

GESTURE=/usr/local/lib/voodooinputmavericks/VoodooInputMavericks.kext
PLIST=/Library/LaunchDaemons/dev.modernmavericks.voodooinputmavericks.plist
STATE=/var/db/voodooinputmavericks-boot.state

chown -R root:wheel "$GESTURE"
chmod -R 755 "$GESTURE"
chmod +x /usr/local/sbin/voodooinputmavericks-run

mkdir -p /var/db
echo ok > "$STATE"

launchctl unload "$PLIST" 2>/dev/null || true
launchctl load -w "$PLIST"

# The loader is WatchPaths-only (no RunAtLoad -- that armed launchd's ~8s boot throttle). `load -w`
# above only REGISTERS the job; a WatchPaths job does not run at load. Touch the trigger to fire it now
# so the kext loads at install/update time (the same path the per-user LaunchAgent uses at login).
# voodooinputmavericks-run's /dev/console guard still no-ops it at the login window (no session) --
# correct: it then loads at the user's next login. Without this, an update leaves the OLD kext torn
# down (preinstall) and the NEW kext unloaded until the next login.
TRIGGER=/usr/local/var/voodooinputmavericks/session.trigger
mkdir -p "$(dirname "$TRIGGER")"
touch "$TRIGGER"

# Post-update: BOUNCE the BT MT2 (close+reopen) so it resumes without a tap. A plain wake no-ops here
# because the kext reload leaves the baseband link UP but the multitouch stream stalled
# (openConnection -> kBluetoothHCIErrorACLConnectionAlreadyExists 0x0b); the bounce forces a full
# re-establish + re-enable. Best-effort (paired CoD-0x594 only, no-ops otherwise). Backgrounded +
# settle so a powered-off device's page timeout can't stall the installer.
( sleep 3; /usr/local/libexec/mt2_linkstated --bounce-once >/dev/null 2>&1 ) &

# USB->BT handoff daemon: wakes deep-idle Bluetooth on cable unplug (no click). Root LaunchDaemon,
# runs regardless of GUI session. Best-effort — a UX feature, never fail the install on it.
HANDOFF=/Library/LaunchDaemons/dev.modernmavericks.voodooinputmavericks.linkstated.plist
chmod +x /usr/local/libexec/mt2_linkstated 2>/dev/null || true
launchctl unload "$HANDOFF" 2>/dev/null || true
launchctl load -w "$HANDOFF" 2>/dev/null || true

# Prefpane live-refresh: shipped ONLY as a SIMBL plugin (we no longer carry a standalone osax loader).
# Drop the bundle into the system SIMBL Plugins dir; a present SIMBLAgent injects it into System
# Preferences. The trackpad itself works without SIMBL — only the Trackpad pane (settings/battery/rename)
# needs it — so if SIMBL is absent we install the plugin anyway and point the user at the Mavericks SIMBL
# build. Best-effort; never fail the install on the prefpane feature.
SIMBL_AGENT="/Library/Application Support/SIMBL/SIMBLAgent.app"
SIMBL_PLUGINS="/Library/Application Support/SIMBL/Plugins"
STAGED_BUNDLE=/usr/local/share/voodooinputmavericks/VoodooInputMavericksPane.bundle
CONSOLE_USER=$(stat -f%Su /dev/console 2>/dev/null)
mkdir -p "$SIMBL_PLUGINS"
rm -rf "$SIMBL_PLUGINS/VoodooInputMavericksPane.bundle"
cp -R "$STAGED_BUNDLE" "$SIMBL_PLUGINS/" 2>/dev/null || true
if [ ! -d "$SIMBL_AGENT" ]; then
    echo "VoodooInputMavericks: the trackpad works now. For the Trackpad preference pane (settings, battery,"
    echo "rename), install SIMBL, then relaunch System Preferences:"
    echo "  https://mavericksforever.com/downloads/SIMBL.pkg"
else
    echo "VoodooInputMavericks: prefpane companion installed as a SIMBL plugin."
fi

# Scheduled background update check (per-user LaunchAgent). Independent of the prefpane route above.
# Both the plist and this load logic are rendered by mavericks-shared-cmake's stage_updater.sh, so all
# the ModernMavericks products load their agent the same way instead of each hand-rolling it. The
# snippet ships beside this script and only exists in a Sparkle build, so its absence is normal and
# means "this build has no updater". It is written to be sourced: no `exit`, MAV_-prefixed variables.
AGENT_LOAD="$(dirname "$0")/agent-load.sh"
if [ -f "$AGENT_LOAD" ]; then
    . "$AGENT_LOAD"
fi

# Connect/disconnect bezel: our BezelServices plugin (installed above) reuses Apple's OWN trackpad art,
# copied here from the co-resident AppleBluetoothMultitouch.plugin — we ship no art (nothing Apple-derived
# in our payload; we're a compatibility layer presenting the MT2 as a genuine Apple trackpad). Best-effort;
# the OSD degrades to image-less if Apple's art is absent.
BEZEL_SRC="/Library/Application Support/Apple/BezelServices/AppleBluetoothMultitouch.plugin/Contents/Resources"
BEZEL_DST="/Library/Application Support/Apple/BezelServices/MavericksMultitouch.plugin/Contents/Resources"
mkdir -p "$BEZEL_DST"
cp "$BEZEL_SRC/BtTrackpad.pdf"     "$BEZEL_DST/MavericksTrackpad.pdf"  2>/dev/null || true
cp "$BEZEL_SRC/MagicTrackpad.icns" "$BEZEL_DST/MavericksTrackpad.icns" 2>/dev/null || true

# Updater icon: same rule, same source. We ship a neutral placeholder and borrow Apple's Magic
# Trackpad art from the machine at install time, so no Apple artwork travels in our payload. The
# updater is what a user sees when an update is offered, so it should look like the device it drives.
UPD_ICON="/Library/Application Support/ModernMavericks/Trackpad2Updater.app/Contents/Resources/Trackpad2Updater.icns"
[ -f "$UPD_ICON" ] && cp "$BEZEL_SRC/MagicTrackpad.icns" "$UPD_ICON" 2>/dev/null || true

exit 0
