# dist/scripts/postinstall
#!/bin/sh
# Activate: fix kext ownership (kextload requires root:wheel) and start the
# daemon. The launchd wrapper (mt2d-run) loads the kext from /usr/local,
# re-enumerates the device if needed, and runs mt2d. No /Library/Extensions,
# no kextcache: the unsigned kext only loads from a loose path.
set -e

KEXT=/usr/local/lib/mt2d/MT2Claim.kext
PLIST=/Library/LaunchDaemons/com.schmonz.mt2d.plist

chown -R root:wheel "$KEXT"
chmod -R 755 "$KEXT"
chmod +x /usr/local/sbin/mt2d-run

launchctl unload "$PLIST" 2>/dev/null || true
launchctl load -w "$PLIST"

exit 0
