# dist/scripts/postinstall
#!/bin/sh
# Activate: fix kext ownership (kextload requires root:wheel), reset the boot
# sentinel to ok (a fresh install should attempt full-gesture mode), and
# start the daemon. The launchd wrapper (voodooinputmavericks-run) loads the one kext from
# /usr/local and re-enumerates the device so the in-kernel readers win their
# matches. No /Library/Extensions, no kextcache: the unsigned kext only loads
# from a loose path.
set -e

GESTURE=/usr/local/lib/voodooinputmavericks/VoodooInputMavericks.kext
PLIST=/Library/LaunchDaemons/dev.modernmavericks.voodooinputmavericks.plist
STATE=/var/db/voodooinputmavericks-boot.state

chown -R root:wheel "$GESTURE"
chmod -R 755 "$GESTURE"
chmod +x /usr/local/sbin/voodooinputmavericks-run

mkdir -p /var/db
echo ok > "$STATE"

launchctl unload "$PLIST" 2>/dev/null || true
launchctl load -w "$PLIST"

# Post-update: BOUNCE the BT MT2 (close+reopen) so it resumes without a tap. A plain wake no-ops here
# because the kext reload leaves the baseband link UP but the multitouch stream stalled
# (openConnection -> kBluetoothHCIErrorACLConnectionAlreadyExists 0x0b); the bounce forces a full
# re-establish + re-enable. Best-effort (paired CoD-0x594 only, no-ops otherwise). Backgrounded +
# settle so a powered-off device's page timeout can't stall the installer.
( sleep 3; /usr/local/libexec/mt2_linkstated --bounce-once >/dev/null 2>&1 ) &

# USB->BT handoff daemon: wakes deep-idle Bluetooth on cable unplug (no click). Root LaunchDaemon,
# runs regardless of GUI session. Best-effort — a UX feature, never fail the install on it.
HANDOFF=/Library/LaunchDaemons/dev.modernmavericks.voodooinputmavericks.linkstated.plist
chmod +x /usr/local/libexec/mt2_linkstated 2>/dev/null || true
launchctl unload "$HANDOFF" 2>/dev/null || true
launchctl load -w "$HANDOFF" 2>/dev/null || true

# Prefpane live-refresh: shipped ONLY as a SIMBL plugin (we no longer carry a standalone osax loader).
# Drop the bundle into the system SIMBL Plugins dir; a present SIMBLAgent injects it into System
# Preferences. The trackpad itself works without SIMBL — only the Trackpad pane (settings/battery/rename)
# needs it — so if SIMBL is absent we install the plugin anyway and point the user at the Mavericks SIMBL
# build. Best-effort; never fail the install on the prefpane feature.
SIMBL_AGENT="/Library/Application Support/SIMBL/SIMBLAgent.app"
SIMBL_PLUGINS="/Library/Application Support/SIMBL/Plugins"
STAGED_BUNDLE=/usr/local/share/voodooinputmavericks/VoodooInputMavericksPane.bundle
CONSOLE_USER=$(stat -f%Su /dev/console 2>/dev/null)
mkdir -p "$SIMBL_PLUGINS"
rm -rf "$SIMBL_PLUGINS/VoodooInputMavericksPane.bundle"
cp -R "$STAGED_BUNDLE" "$SIMBL_PLUGINS/" 2>/dev/null || true
if [ ! -d "$SIMBL_AGENT" ]; then
    echo "VoodooInputMavericks: the trackpad works now. For the Trackpad preference pane (settings, battery,"
    echo "rename), install SIMBL, then relaunch System Preferences:"
    echo "  https://mavericksforever.com/downloads/SIMBL.pkg"
else
    echo "VoodooInputMavericks: prefpane companion installed as a SIMBL plugin."
fi

# Scheduled background update check (per-user LaunchAgent). Independent of the prefpane route above.
# Only meaningful when the updater app was staged (Sparkle-enabled build); load it for the current GUI
# user so the first daily check can run without a logout. Best-effort — never fail the install on it.
UPDATE_AGENT=/Library/LaunchAgents/dev.modernmavericks.voodooinputmavericks.updatecheck.plist
UPDATER_APP=/usr/local/lib/voodooinputmavericks/MavericksTrackpad2Updater.app/Contents/MacOS/MavericksTrackpad2Updater
if [ -x "$UPDATER_APP" ] && [ -n "$CONSOLE_USER" ] && [ "$CONSOLE_USER" != "root" ]; then
    sudo -u "$CONSOLE_USER" launchctl unload "$UPDATE_AGENT" 2>/dev/null || true
    sudo -u "$CONSOLE_USER" launchctl load -w "$UPDATE_AGENT" 2>/dev/null || true
fi

exit 0
