# dist/com.schmonz.mt2d.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.schmonz.mt2d</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/mt2d-run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!-- The wrapper loads the kext, confirms health, and exits; the kext runs
         in-kernel on its own. Run once at boot, do NOT respawn. -->
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/var/log/mt2d.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/mt2d.log</string>
</dict>
</plist>
