# dist/com.schmonz.mt2panewatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Per-user GUI-session agent: watches for System Preferences launches and sends
         our MT2x/load Apple event so the prefpane-refresh osax loads (no SIMBL). Must be
         a LaunchAgent (not Daemon) so it runs in the Aqua session where NSWorkspace
         launch notifications are delivered. -->
    <key>Label</key>
    <string>com.schmonz.mt2panewatch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/libexec/mt2_pane_watch</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
