# dist/com.schmonz.mt2d.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.schmonz.mt2d</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/mt2d-run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!-- Runs at boot (RunAtLoad -> a no-op behind mt2d-run's /dev/console session guard, so the login
         screen stays on Apple's generic HID) AND when the per-user LaunchAgent touches this trigger at
         login (session start -> the guard passes -> full multitouch takeover). -->
    <key>WatchPaths</key>
    <array>
        <string>/usr/local/var/mt2d/session.trigger</string>
    </array>
    <!-- The wrapper loads the kext, confirms health, and exits; the kext runs
         in-kernel on its own. Do NOT respawn. -->
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/var/log/mt2d.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/mt2d.log</string>
</dict>
</plist>
