# dist/com.schmonz.voodooinputmavericks.session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Per-user Aqua LaunchAgent: at login it touches the trigger the root voodooinputmavericks daemon WatchPaths,
         so voodooinputmavericks-run runs post-session and takes the MT2 over for multitouch. Pre-login (no Aqua
         session) it never runs, so the login screen stays on Apple's generic HID. -->
    <key>Label</key>
    <string>com.schmonz.voodooinputmavericks.session</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>RunAtLoad</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/touch</string>
        <string>/usr/local/var/voodooinputmavericks/session.trigger</string>
    </array>
</dict>
</plist>
