# dist/dev.modernmavericks.voodooinputmavericks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.modernmavericks.voodooinputmavericks</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/voodooinputmavericks-run</string>
    </array>
    <!-- No RunAtLoad, deliberately. A boot-time run would only hit voodooinputmavericks-run's
         /dev/console session guard and no-op out (login screen stays on Apple's generic HID) — but that
         wasted launch armed launchd's ~10s respawn throttle, so the real login-triggered load was then
         delayed ~8s (measured boot->gestures 2026-07-23: ~15s, ~8s of it pure throttle). The job is
         purely WatchPaths-driven: the per-user LaunchAgent touches this trigger at login (session start
         -> the guard passes -> full multitouch takeover), which is the only path that does real work. -->
    <key>WatchPaths</key>
    <array>
        <string>/usr/local/var/voodooinputmavericks/session.trigger</string>
    </array>
    <!-- The wrapper loads the kext, confirms health, and exits; the kext runs
         in-kernel on its own. Do NOT respawn. -->
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/var/log/voodooinputmavericks.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/voodooinputmavericks.log</string>
</dict>
</plist>
